# escape-room-display/escape-room-display.service
[Unit]
Description=Escape room display HTTP (copperdragons slave)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cde
WorkingDirectory=/home/cde/escape-room-display
# Optional: load env vars from a file
# EnvironmentFile=-/home/cde/escape-room-display/.env
# Or set a few inline:
# Environment=DISPLAY_API_KEY=your-secret
# Environment=USE_SUDO=1
ExecStart=/home/cde/escape-room-display/.venv/bin/uvicorn server:app --host 0.0.0.0 --port 8765
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
